# deploy/firedancer-health-exporter.service
[Unit]
Description=Firedancer Health Prometheus Exporter
Documentation=https://github.com/antipov/firedancer-health-exporter
After=network.target firedancer.service
Wants=firedancer.service

[Service]
Type=simple
User=firedancer
Group=systemd-journal

ExecStart=/opt/firedancer-health-exporter/bin/firedancer-exporter \
    --port 9100 \
    --interval 60 \
    --enable-rpc-metrics \
    --rpc-url https://api.mainnet-beta.solana.com \
    --vote-account <YOUR_VOTE_ACCOUNT_PUBKEY> \
    --identity    <YOUR_IDENTITY_PUBKEY>

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=firedancer-exporter

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadOnlyPaths=/
BindReadOnlyPaths=/run/log/journal /var/log/journal

[Install]
WantedBy=multi-user.target
